# TP3-Sécu/aide.service
[Unit]
Description=Exécuter un check AIDE
Wants=aide.timer

[Service]
Type=simple
ExecStart=/usr/bin/aide --check